# 1.1/80-test.rules
KERNEL=="sdb", RUN+="/bin/device_added.sh"